# data/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx
redirect_stderr=true

[program:uwsgi]
command =/usr/local/bin/uwsgi --ini  /var/www/flask/uwsgi.ini
redirect_stderr=true

[program:flaskdebug]
command=/usr/bin/python3 /var/www/flask/debug.py
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
startretries=100
